# templates/rh_initscript.erb
#! /bin/sh
#
# simple vpnc init.d script for rh systems
#
# chkconfig: 2345 10 90
# description: Starts and stops VPNC service

# Source function library.
. /etc/init.d/functions

# Check that networking is up.
[ "${NETWORKING}" = "no" ] && exit 6

case "$1" in
  start)
        echo -n "Starting VPN connection" 
        /usr/sbin/vpnc /etc/vpnc/<%= @name %>.conf
        exit $?
    ;;
  stop)
        echo -n "Stopping VPN connection" 
        /usr/sbin/vpnc-disconnect
        exit $?
    ;;
  *)
    echo "Usage: /etc/init.d/vpnc-<%= @name %> {start|stop}"
    exit 1
    ;;
esac

exit 0

